# deploy/wecom-webhook.service
[Unit]
Description=WeCom Bot Webhook
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sg audio -c \"/home/jxy001a1/code/knowledge-wiki/.venv/bin/python3 -m knowledge_wiki webhook\"
WorkingDirectory=/home/jxy001a1/code/knowledge-wiki
Restart=always
RestartSec=5
Environment="PYTHONPATH=/home/jxy001a1/code/knowledge-wiki/src"
EnvironmentFile=/home/jxy001a1/code/knowledge-wiki/.env

[Install]
WantedBy=default.target
